# supervisord.conf
[supervisord]
nodaemon=true

[program:zoneminder]
command=service zoneminder start

[program:apache2]
command=service apache2 start